# send_notes/send_notes.service
[Unit]
Description=HiveTalk Send Notes Service
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/scheduler/send_notes
ExecStart=/root/scheduler/send_notes/run.sh
Restart=always
RestartSec=30
Environment=PAGER=cat

# Hardening measures - temporarily disabled for troubleshooting
# PrivateTmp=true
# ProtectSystem=full
# NoNewPrivileges=true
# ProtectHome=true

[Install]
WantedBy=multi-user.target
